# bbb-config/data/etc/systemd/system/bbb-transcode-akka.service.d/override.conf
[Unit]
Requires=redis-server.service
After=redis-server.service
